# app/mywebapp.service
[Unit]
Description=My Web App (Node.js)
After=network-online.target
Wants=network-online.target

[Service]
Environment=PORT=3000
Type=simple
ExecStart=/usr/bin/node /opt/mywebapp/server.js
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
